# GameContainer/app/supervisord/nginx.proxy.conf
[program:nginx]
command=/usr/local/nginx/sbin/nginx
# will be started once the nginx configs have been copied
autostart=false
autorestart=true
