# GameContainer/app/supervisord/setup.conf
[program:setup]
command=/opt/app/setup.sh
autostart=true
autorestart=unexpected
